# roles/augmenture/templates/augmenture.service.j2
[Unit]
Description=Augmenture webapp

[Service]
Type=simple
ExecStart=/usr/bin/augmenture
ExecStop= /bin/kill -2 $MAINPID
